# ops/quant-monitor/systemd/codex-daily-briefing.service.example
[Unit]
Description=QuantStrategyLab daily strategy briefing (task 10)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Projects/AIAuditBridge/ops/quant-monitor
Environment=QUANT_MONITOR_ROOT=/home/ubuntu/Projects/AIAuditBridge/ops/quant-monitor
Environment=AIAUDIT_BRIDGE_ROOT=/home/ubuntu/Projects/AIAuditBridge
Environment=QUANT_PLATFORM_KIT_ROOT=/home/ubuntu/Projects/QuantPlatformKit
Environment=QUANT_SENTINEL_TELEGRAM_SECRET_NAME=quant-sentinel-telegram-bot-token
Environment=QUANT_SENTINEL_GCP_PROJECT=firstradequant
Environment=GLOBAL_TELEGRAM_CHAT_ID=
Environment=QSL_GITHUB_REPO=QuantStrategyLab/CnEquityStrategies
RuntimeDirectory=quant-monitor
RuntimeDirectoryMode=0750
ExecStartPre=/bin/bash /home/ubuntu/Projects/AIAuditBridge/ops/quant-monitor/scripts/load_telegram_env.sh /run/quant-monitor/telegram.env
EnvironmentFile=-/run/quant-monitor/telegram.env
ExecStart=/bin/bash /home/ubuntu/Projects/AIAuditBridge/ops/quant-monitor/scripts/daily_briefing_pipeline.sh

[Install]
WantedBy=multi-user.target
