# ops/quant-monitor/systemd/codex-quant.timer.example
[Unit]
Description=Run quant monitor every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Unit=codex-quant.service

[Install]
WantedBy=timers.target
