# ops/codex-audit/systemd/model-catalog-sync.service.example
[Unit]
Description=Monthly auto-sync of AI model catalog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/codex-audit-bridge
StateDirectory=codex-audit-bridge
Environment=MODEL_CATALOG_PATH=/var/lib/codex-audit-bridge/model_catalog.json
# Provide OPENAI_API_KEY / ANTHROPIC_API_KEY via this env file (mode 0600).
EnvironmentFile=-/etc/codex-audit-bridge/model-catalog.env
ExecStart=/opt/codex-audit-bridge/.venv/bin/python scripts/sync_model_catalog.py --force
User=codex-audit
Group=codex-audit

[Install]
WantedBy=multi-user.target
