# ops/codex-audit/systemd/model-catalog-sync.service.example
[Unit]
Description=Monthly auto-sync of AI model catalog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Projects/AIAuditBridge
StateDirectory=codex-audit-bridge
Environment=MODEL_CATALOG_PATH=/var/lib/codex-audit-bridge/model_catalog.json
Environment=PYTHONPATH=/home/ubuntu/Projects/AIAuditBridge
# Provider keys are written by deploy_model_catalog_sync.sh from GitHub secrets.
EnvironmentFile=-/etc/codex-audit-bridge/model-catalog.env
ExecStart=/usr/bin/python3 /home/ubuntu/Projects/AIAuditBridge/scripts/sync_model_catalog.py --force
Nice=10

[Install]
WantedBy=multi-user.target
